# root2/etc/s6-overlay/s6-rc.d/init-config-nfs/run
#!/usr/bin/with-contenv bash

# NFS mounts
mount -t nfs -o nolock $NFSPATH /media